# bin/daemon
#! /bin/sh

### BEGIN INIT INFO
# Provides: ant
# Required-Start: $remote_fs $syslog
# Required-Stop:  $remote_fs $syslog
# Default-Start:  2 3 4 5
# Default-Stop:
# Short-Description:  start ant daemon
### END INIT INFO

RUBY_HOME=/opt/ruby

PID=/var/run/ruby_daemon.pid
DAEMON=$RUBY_HOME/bin/ant

test -x $DAEMON || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
    if init_is_upstart; then
      exit 1
    fi

    log_daemon_msg "Starting ant daemon" ant_daemon || true

    if start-stop-daemon --start --quiet --oknodo --background --pidfile $PID --exec $DAEMON -- daemon; then
      log_end_msg 0 || true
    else
      log_end_msg 1 || true
    fi

    ;;
  stop)
    if init_is_upstart; then
      exit 0
    fi

    log_daemon_msg "Stopping ant daemon" ant_daemon || true

    if start-stop-daemon --stop --quiet --oknodo --pidfile $PID; then
      log_end_msg 0 || true
    else
      log_end_msg 1 || true
    fi

    ;;
  restart)
    if init_is_upstart; then
      exit 1
    fi

    $0 stop
    sleep 1
    $0 start

    ;;
  status)
    status_of_proc -p $PID $DAEMON ant_daemon
    exit $?

    ;;
  *)
    log_action_msg "Usage: /etc/init.d/ant_daemon {start|stop|restart|status}" || true
    exit 1
esac

exit 0