# bin/daemon.service
[Unit]
Description=ant daemon
After=remote-fs.target syslog.target

[Service]
TimeoutStartSec=60
ExecStart=${RUBY_HOME}/bin/ruby ${RUBY_HOME}/bin/ant daemon

[Install]
WantedBy=multi-user.target